# supervisord.conf
[supervisord]
nodaemon=true

[program:nodejs]
directory=/tmp/myapp
command=nodejs server.js
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/nodejs.log
