# supervisord.conf
[supervisord]
nodaemon=true

[program:nodejs]
directory=/tmp/myapp
command=nodejs server.js
autostart=true
autorestart=true
stderr_logfile=/var/log/nodejs.log
